# .config/systemd/user/dunst.service.d/override.conf
[Unit]
After=niri.service
StartLimitIntervalSec=0

[Service]
Restart=on-failure
RestartSec=2s
